# debian/ntpsec.ntp-rotate-stats.service
[Unit]
Description=Rotate ntpd stats
Requisite=ntp.service

[Service]
Type=simple
ExecStart=/usr/lib/ntp/rotate-stats
