# debian/ntpsec.cron.d
25 6     * * * root [ ! -d /run/systemd/system ] && [ -x /usr/lib/ntp/rotate-stats ] && /usr/lib/ntp/rotate-stats
